# files/etc/systemd/system/k3s.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/k3s server \
            --disable=servicelb,traefik \
            --tls-san=routerd.erdii.net
